# cron.d/check-resolv
*/30 * * * * root /usr/local/bin/check-resolv
