# cron.d/warn_battery_low
*/15 * * * * ###USER### /usr/local/bin/warn_battery_low
