# cron.d/warn-battery-low
*/15 * * * * ###USER### /usr/local/bin/warn-battery-low
